# .archive/worktrees/agent-1.3/aurigraph-av10-7/dist-dev4/aurigraph-dev4.service
[Unit]
Description=Aurigraph V10 Dev4 - dlt.aurigraph.io
After=network.target

[Service]
Type=simple
User=aurigraph
WorkingDirectory=/var/www/dlt.aurigraph.io
Environment="NODE_ENV=production"
Environment="PORT=8080"
Environment="HOST=0.0.0.0"
ExecStart=/usr/bin/node dist/index-classical-simple.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
